# ansible/roles/forge-neo/templates/forge-neo.service.j2
[Unit]
Description=Stable Diffusion WebUI Forge Neo
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ forge_install_dir }}/app
Environment="VIRTUAL_ENV={{ forge_install_dir }}/app/venv"
Environment="PATH={{ forge_install_dir }}/app/venv/bin:{{ ansible_env.HOME }}/.cargo/bin:/usr/local/cuda-12.8/bin:{{ ansible_env.PATH }}"
Environment="LD_LIBRARY_PATH=/usr/local/cuda-12.8/lib64"
ExecStart={{ forge_install_dir }}/app/venv/bin/python3 launch.py --listen --port 7860
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
