# ansible/roles/forge-neo/templates/forge-neo.service.j2
[Unit]
Description=AUTOMATIC1111 Stable Diffusion WebUI
After=network.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ forge_install_dir }}/app
Environment="VIRTUAL_ENV={{ forge_install_dir }}/app/venv"
Environment="PATH={{ forge_install_dir }}/app/venv/bin:{{ ansible_env.HOME }}/.cargo/bin:{{ ansible_env.PATH }}"
Environment="PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:512"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONFAULTHANDLER=1"
ExecStart=/bin/bash {{ forge_install_dir }}/app/start-forge.sh
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-neo

[Install]
WantedBy=multi-user.target
